# templates/etc/update-motd.d/20-services.j2
#!/bin/bash

{% if _motd.services is defined %}
Services
{% for name in _motd.services.services %}
{{ _motd.services.command|replace("@@NAME@@", name) }}
{% endfor %}
echo
{% endif %}

